# cansetup.conf
# can4linux cansetup starter
start on local-filesystems or runlevel [2345]
stop on runlevel [!2345]

task

script
	exec /usr/sbin/cansetup
end script
